# etc/supervisor.d/02-memcached.conf
[program:memcached]
command=/usr/bin/memcached -u memcache -l 127.0.0.1 -m 64 -p 11211
stdout_events_enabled=true
stderr_events_enabled=true
